# root/etc/services.d/apache/run
#!/usr/bin/with-contenv bash

#source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND
apachectl start

