# root/etc/services.d/sshd/run
#!/usr/bin/with-contenv bash
mkdir -p /var/run/sshd
/usr/sbin/sshd